# archivos-raspi/autostart/pd.desktop
[Desktop Entry]
Type=Application
Name=PureData
Exec=bash -c "sleep 10s; pd -rt -jack /home/vcnt/pi/Londres38/_MAIN.pd"

