# archivos-raspi/autostart/qjackctl.desktop
[Desktop Entry]
Type=Application
Name=qjackctl
Exec=bash -c "qjackctl &"
